# .devcontainer/supervisord.conf
[program:desktop]
command=/usr/local/bin/start-desktop.sh
autostart=true
autorestart=true
user=vscode
redirect_stderr=true
stdout_logfile=/home/vscode/desktop.log
environment=DISPLAY=":1",HOME="/home/vscode",VNC_PASSWORD="${VNC_PASSWORD}"
